# files/acg-httpd
#!/sbin/openrc-run

acg_path="/opt/acg"

# Load Configuration
if [ -f "${acg_path}/files/acg-cfg" ]; then
  . "${acg_path}/files/acg-cfg"
else
  echo "The configuration file is missing. Please re-run the installation script."
  exit 1
fi
# Load Version
. "${acg_path}/files/version"

command="/usr/sbin/httpd"
command_args="-p ${ACG_API_PORT} -h ${acg_path}/www"
name="Alpine Clash Gateway REST API"
description="HTTP Server for Alpine Clash Gateway REST API"

depend() {
	need net localmount
}
